# root/etc/service/opensim-robust-server/run
#!/bin/sh
#
# yes, sleeps are inefficient and potentially inaccurate.  but better
# than nothing for now.

SERVERNAME=opensim-robust-server
SERVER=OpenSim.Server.exe
SERVEROPTS=
OPENSIMDIR=/usr/lib/opensim

if `svok /etc/service/opensim-grid-server`; then

    sleep 10
    echo [ `date` ] Starting $SERVERNAME ...
    su - opensim -c "cd ${OPENSIMDIR}/bin; screen -S $SERVERNAME -t $SERVERNAME -D -m mono $SERVER $SERVEROPTS"
    echo [ `date` ] Restarting $SERVERNAME ...

else

    sleep 2

fi
